# roles/claude-code/templates/claude-code.container.j2
[Unit]
Description=Claude Code Container
After=network-online.target
Wants=network-online.target

[Container]
Image={{ claude_code_image }}:{{ claude_code_version }}
ContainerName=claude-code
User={{ claude_code_user }}

# Workspace persists across sessions; no host filesystem access
Volume={{ claude_code_workspace_volume }}:/home/claude/workspace:z

# API key from env file (root-only, mode 0600)
EnvironmentFile={{ claude_code_env_file }}

# Security hardening
NoNewPrivileges=true
DropCapability=all

# tmux keeps the container alive; attach with: podman attach claude-code
# Then run: claude --dangerously-skip-permissions
PodmanArgs=--tty

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target default.target
