# roles/claude-code/templates/claude-code.container.j2
[Unit]
Description=Claude Code Container
After=network-online.target
Wants=network-online.target

[Container]
Image={{ claude_code_image }}:{{ claude_code_version }}
ContainerName=claude-code
User={{ claude_code_user }}

# Workspace and auth config persist across container recreations
Volume={{ claude_code_workspace_volume }}:/home/claude/workspace:z
Volume=claude-code-config:/home/claude/.claude:z

# API key from env file (root-only, mode 0600)
EnvironmentFile={{ claude_code_env_file }}

# Security hardening
NoNewPrivileges=true
DropCapability=all

# tmux keeps the container alive; attach with: podman attach claude-code
# Then run: claude --dangerously-skip-permissions
Environment=TERM=xterm-256color
Environment=SHELL=/bin/bash
Environment=LANG=en_US.UTF-8
Environment=LC_ALL=en_US.UTF-8
PodmanArgs=--tty

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target default.target
